# services_inviting.service
[Unit]
Description=services inviting

[Service]
User=q
Type=simple
WorkingDirectory=/home/q/services/services-inviting
ExecStart=/usr/bin/python3 /home/q/services/services-inviting/main.py
Restart=always

[Install]
WantedBy=multi-user.target
